# sv/syslogd/run
#!/bin/sh
exec chpst -Unobody syslogd -n -O -
