# sv/syslog-ng/run
#!/bin/sh

exec syslog-ng --stderr --foreground
