# sv/socklog-unix/run
#!/bin/sh
exec chpst -Unobody socklog unix /dev/log
